# services/netdiscover/run
#!/usr/bin/with-contenv /bin/bash
test -f /env && source /env
test -f /env && export $(cut -d= -f1 /env)
echo 'Starting netdiscover'
exec python3 /netdiscover.py



