# services/netdiscover/run
#!/usr/bin/with-contenv /bin/bash
echo 'Starting netdiscover'
exec python3 /netdiscover.py



